# templates/node_exporter.service.j2
[Unit]
Description=NodeExporter

[Service]
TimeoutStartSec=0
ExecStart={{ dest_path}}/node_exporter

[Install]
WantedBy=multi-user.target